# api/crontab
*/5 * * * * root cd /api && ./updatePagesStatus.py > listener.log 2>&1

