# api/crontab
*/5 * * * * root cd /api && ./updatePagesStatus.py > updater.log 2>&1

